# duplicati/rootfs/etc/services.d/duplicati/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Duplicati
# ==============================================================================


bashio::log.info "Starting Duplicati..."
cd /app/duplicati || exit
exec mono Duplicati.Server.exe \
	--webservice-interface=any --server-datafolder=/data --webservice-allowed-hostnames=*