# duplicati/rootfs/etc/services.d/duplicati/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Duplicati
# ==============================================================================


bashio::log.info "Starting Duplicati..."

#exec php-fpm7 --nodaemonize -R